# run/var-mnt-server.automount.j2
[Unit]
Description=Server storage automount

[Automount]
Where=/var/mnt/server
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
